# backend/blog-backend.service
[Unit]
Description=Azhen Blog Backend Service
After=network.target mysql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/azhen-blog/backend
ExecStart=/usr/bin/java -jar /var/www/azhen-blog/backend/target/blog-backend-1.0.0.jar
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=azhen-blog-backend

# 环境变量
Environment="SPRING_DATASOURCE_URL=jdbc:mysql://localhost:3306/azhen_blog?useSSL=false&serverTimezone=UTC&allowPublicKeyRetrieval=true"
Environment="SPRING_DATASOURCE_USERNAME=root"
Environment="SPRING_DATASOURCE_PASSWORD=yourpassword"
Environment="SERVER_PORT=8080"

[Install]
WantedBy=multi-user.target
